# init/systemd/bump_instance.timer
[Unit]
Description=Run bump_instance every minute

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
