# init/systemd/bump_instance.timer
[Unit]
Description=Run bump_instance every minute

[Timer]
# run every 3 hours
OnCalendar=*-*-* *:00/3:00:00

[Install]
WantedBy=timers.target
