# init/systemd/bump_instance.timer
[Unit]
Description=bump_instance timer

[Timer]
OnCalendar=*-*-* 0,3,6,9,12,15,18,21:00:00

[Install]
WantedBy=timers.target
