# init/systemd/feed_all_from_db.timer
[Unit]
Description=Run every minute

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
