# init/systemd/feed_bulk_from_db.timer
[Unit]
Description=Run feed_bulk_from_db every minute

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
